# deploy/provisioning/systemd/cyprine-backend.service
[Unit]
Description=Cyprine Heroes FastAPI backend
After=network.target network-online.target
Wants=network-online.target
StartLimitIntervalSec=300
StartLimitBurst=5

[Service]
Type=simple
User=cyprine
Group=cyprine
WorkingDirectory=/opt/cyprine-heroes/backend
EnvironmentFile=/etc/cyprine-heroes/backend.env

# Pre-flight checks
ExecStartPre=/opt/cyprine-heroes/deploy/scripts/validate-env.sh
ExecStartPre=/opt/cyprine-heroes/deploy/scripts/pre-start-check.sh

# Main service
ExecStart=/opt/cyprine-heroes/venv/bin/uvicorn app.main:app --host 127.0.0.1 --port 8000 --workers 2

# Health check
ExecStartPost=/bin/sleep 5
ExecStartPost=/opt/cyprine-heroes/deploy/scripts/health-check.sh

# Graceful shutdown
ExecStop=/bin/kill -TERM $MAINPID
KillMode=mixed
TimeoutStopSec=30

# Restart policy
Restart=always
RestartSec=10

# Security hardening
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/opt/cyprine-heroes/backend/uploads /tmp
RestrictSUIDSGID=yes
LockPersonality=yes
MemoryDenyWriteExecute=yes

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=cyprine-backend

# Resource limits
LimitNOFILE=65536
MemoryMax=1G

[Install]
WantedBy=multi-user.target
